# apps/router-plugin/root/etc/init.d/nextpanel-sync
#!/bin/sh /etc/rc.common
# NextPanel HomeProxy node sync service
# Writes the NextPanel subscription URL into HomeProxy's UCI subscription list
# and triggers HomeProxy's built-in node import.

START=99
STOP=01
EXTRA_COMMANDS="sync force_sync"
EXTRA_HELP="	sync		Sync nodes from panel (respects enabled flag)
	force_sync	Force sync regardless of enabled setting"

LOGFILE=/var/log/nextpanel-sync.log
LOGMAX=200
CRON_FILE=/etc/crontabs/root
CRON_TAG="nextpanel-sync"
HP_SCRIPTS="/etc/homeproxy/scripts"

log() {
    local msg="$(date '+%Y-%m-%d %H:%M:%S') [nextpanel] $1"
    echo "$msg" >> "$LOGFILE"
    local lines
    lines=$(wc -l < "$LOGFILE" 2>/dev/null || echo 0)
    if [ "$lines" -gt "$LOGMAX" ]; then
        local tmp="/tmp/nextpanel-log-trim-$$.tmp"
        tail -n "$LOGMAX" "$LOGFILE" > "$tmp" && mv "$tmp" "$LOGFILE"
    fi
}

get_cfg() {
    uci -q get "nextpanel.config.$1"
}

set_cfg() {
    uci -q set "nextpanel.config.$1=$2"
    uci -q commit nextpanel
}

# ─── Core sync logic ──────────────────────────────────────────────────────────

do_sync() {
    local enabled
    enabled=$(get_cfg enabled)
    local url
    url=$(get_cfg subscription_url)

    if [ "$1" != "force" ] && [ "$enabled" != "1" ]; then
        log "Sync is disabled — skipping"
        return 0
    fi

    if [ -z "$url" ]; then
        log "ERROR: subscription_url is not configured"
        set_cfg last_status "error"
        set_cfg last_error "未配置订阅链接"
        return 1
    fi

    log "Syncing nodes from: $url"

    # Write our subscription URL into HomeProxy's subscription_url list.
    # HomeProxy stores multiple URLs as a list; we own a single entry tagged
    # with a comment marker so we can replace it on each sync without
    # disturbing any other URLs the user may have added.
    local hp_cfg="homeproxy"

    # Ensure the subscription section exists
    uci -q get "${hp_cfg}.subscription" > /dev/null 2>&1 || {
        uci -q set "${hp_cfg}.subscription=homeproxy"
    }

    # Remove any previous NextPanel-managed URL (identified by our tag in the list)
    local prev
    prev=$(uci -q get "${hp_cfg}.subscription.subscription_url" 2>/dev/null)
    uci -q delete "${hp_cfg}.subscription.subscription_url" 2>/dev/null || true

    # Re-add non-NextPanel URLs (if user had others) then add ours
    # For simplicity: replace the entire list with our URL.
    # If users need multiple subscriptions they can add them in HomeProxy UI after sync.
    uci -q add_list "${hp_cfg}.subscription.subscription_url=${url}"
    uci -q commit "${hp_cfg}"

    log "Updated homeproxy.subscription.subscription_url"

    # Trigger HomeProxy's subscription update script
    if [ -f "${HP_SCRIPTS}/update_subscriptions.uc" ]; then
        log "Running HomeProxy subscription update..."
        if ucode "${HP_SCRIPTS}/update_subscriptions.uc" >> "$LOGFILE" 2>&1; then
            log "Subscription update complete"
            set_cfg last_sync "$(date '+%Y-%m-%d %H:%M:%S')"
            set_cfg last_status "success"
            set_cfg last_error ""
            # Reload HomeProxy to apply new nodes
            if [ -f /etc/init.d/homeproxy ]; then
                log "Reloading HomeProxy..."
                /etc/init.d/homeproxy restart
                log "HomeProxy reloaded"
            fi
        else
            log "ERROR: HomeProxy subscription update failed"
            set_cfg last_status "error"
            set_cfg last_error "HomeProxy 节点导入失败，请查看日志"
            return 1
        fi
    else
        log "ERROR: HomeProxy scripts not found at ${HP_SCRIPTS}"
        set_cfg last_status "error"
        set_cfg last_error "未找到 HomeProxy，请先安装 HomeProxy"
        return 1
    fi
}

# ─── Cron management ─────────────────────────────────────────────────────────

setup_cron() {
    local interval
    interval=$(get_cfg interval)

    sed -i "/${CRON_TAG}/d" "$CRON_FILE" 2>/dev/null

    case "$interval" in
        3600)  echo "0 * * * * /etc/init.d/nextpanel-sync sync # ${CRON_TAG}" >> "$CRON_FILE" ;;
        21600) echo "0 */6 * * * /etc/init.d/nextpanel-sync sync # ${CRON_TAG}" >> "$CRON_FILE" ;;
        86400) echo "0 3 * * * /etc/init.d/nextpanel-sync sync # ${CRON_TAG}" >> "$CRON_FILE" ;;
        *)     log "Auto-refresh disabled (interval=$interval)" ;;
    esac

    if [ -f /var/run/crond.pid ]; then
        /etc/init.d/cron restart 2>/dev/null || true
    fi
}

remove_cron() {
    sed -i "/${CRON_TAG}/d" "$CRON_FILE" 2>/dev/null
    if [ -f /var/run/crond.pid ]; then
        /etc/init.d/cron restart 2>/dev/null || true
    fi
}

# ─── rc.common entry points ───────────────────────────────────────────────────

start() {
    local enabled
    enabled=$(get_cfg enabled)
    [ "$enabled" != "1" ] && { log "Service disabled — not starting"; return 0; }
    log "Service starting"
    setup_cron
    do_sync
}

stop() {
    log "Service stopping"
    remove_cron
}

restart() {
    stop
    start
}

sync() {
    do_sync
}

force_sync() {
    do_sync force
}
